# docker/supervisor.d/mongod.conf
[program:mongod]
autostart=false
user=mongod
command=/usr/bin/mongod -f /etc/mongod.conf
