# docker/profile.d/sudo.sh
# enables alias replacement after `sudo`
alias sudo='sudo '

# repeats the previous command with sudo :)
alias please='sudo $(history -p !!)'
alias pls='please'
